# net-p2p/electrs/pkg-plist
bin/electrs
%%ETCDIR%%.toml.sample
etc/rc.d/electrs
